# scripts/tuitional-app.service
[Unit]
Description=Tuitional AI OS — Backend Server
After=network.target

[Service]
Type=simple
User=tuitional
WorkingDirectory=/home/tuitional/app/backend
EnvironmentFile=/home/tuitional/app/backend/.env.production
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
