# scripts/tuitional-daemon.service
[Unit]
Description=Tuitional AI OS — Task Polling Daemon
After=tuitional-app.service

[Service]
Type=simple
User=tuitional
WorkingDirectory=/home/tuitional/app
EnvironmentFile=/home/tuitional/app/backend/.env.production
ExecStart=/usr/bin/node scripts/task-daemon.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
